# external/rootfs_overlay/common/etc/profile.d/hello.sh
#!/bin/sh
#
# Display version of OS running

. /etc/version

echo "--------------------"
echo "$VERSION ($COMMIT)"
echo "--------------------"

